# cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
Type=simple
User=root
# Replace with the actual path to your cloudflared binary (use 'which cloudflared' to find it)
ExecStart=/usr/local/bin/cloudflared tunnel run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target